# isr/runtime/systemd/inspect-manager.service
[Unit]
Description=Inspect Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cloud
WorkingDirectory=/home/cloud
ExecStart=/home/cloud/workspaces/inspect-slide-rest/runtime/bin/start-inspect-manager.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
